# dns_updater/crontab
@reboot minecraft python3 $MINECRAFT_TOOLS_HOME/dns_updater/dns_updater.py
