# dns_updater/crontab
@reboot minecraft /usr/local/bin/dns_updater.py
