# dns_updater/crontab
@reboot minecraft $MINECRAFT_TOOLS_HOME/dns_updater/dns_updater.py
