# playercount_metric/crontab
* * * * * minecraft /opt/minecraft_aws_tools/playercount_metric/playercount_metric.sh

