# playercount_metric/crontab
* * * * * minecraft python3 $MINECRAFT_HOME/minecraft_aws_tools/playercount_metric/paper_updater.py
