# playercount_metric/crontab
* * * * * minecraft /bin/sh $MINECRAFT_TOOLS_HOME/playercount_metric/playercount_metric.sh

