# paper_updater/crontab
@reboot minecraft python3 $MINECRAFT_TOOLS_HOME/paper_updater/paper_updater.py
