# paper_updater/crontab
@reboot minecraft python3 $MINECRAFT_HOME/minecraft_aws_tools/paper_updater/paper_updater.py
